# scripts/directus-stack.service
[Unit]
Description=Directus Docker Compose Stack
Requires=docker.service
After=docker.service

[Service]
WorkingDirectory=/vagrant
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=5s
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
